# init/systemd/network-dmenu-privileged.service
# network-dmenu systemd user service (privileged version)
# Use this version if the standard service fails with permission errors
# Install to: ~/.config/systemd/user/network-dmenu.service
# Enable with: systemctl --user enable --now network-dmenu

[Unit]
Description=Network Dmenu Geofencing Daemon (Privileged)
Documentation=https://github.com/user/network-dmenu
After=graphical-session.target network-online.target
Wants=network-online.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=HOME=%h
Environment=USER=%i
Environment=XDG_RUNTIME_DIR=%t
Environment=XDG_DATA_HOME=%h/.local/share
Environment=XDG_CONFIG_HOME=%h/.config
Environment=RUST_LOG=debug
ExecStart=%h/.local/bin/network-dmenu --daemon --log-level debug
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Enhanced device and network access permissions
DevicePolicy=auto
SupplementaryGroups=netdev wheel
PrivateNetwork=false
PrivateDevices=false
ProtectKernelTunables=false
ProtectControlGroups=false
RestrictSUIDSGID=false

# Capabilities needed for network operations
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_DAC_OVERRIDE

# File system access
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=%h/.local/share/network-dmenu %h/.config/network-dmenu /tmp
PrivateMounts=false

# Resource limits
MemoryLimit=100M
CPUQuota=10%

[Install]
WantedBy=graphical-session.target