# init/systemd/network-dmenu.service
# network-dmenu systemd user service
# Install to: ~/.config/systemd/user/network-dmenu.service
# Enable with: systemctl --user enable --now network-dmenu

[Unit]
Description=Network Dmenu Geofencing Daemon
Documentation=https://github.com/user/network-dmenu
After=graphical-session.target network.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=HOME=%h
Environment=RUST_LOG=debug
ExecStart=%h/.local/bin/network-dmenu --daemon --log-level debug
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Essential permissions for network device access
PrivateDevices=false
PrivateNetwork=false
SupplementaryGroups=netdev
ReadWritePaths=%h/.local/share/network-dmenu /tmp

# Security settings (minimal for functionality)
NoNewPrivileges=true
ProtectKernelTunables=false
RestrictSUIDSGID=false

[Install]
WantedBy=graphical-session.target